# roles/es-master/files/elasticsearch.service.d/override.conf
[Service]
LimitMEMLOCK=infinity
PermissionsStartOnly=true
